# templates/dev.local.iterm2-gpu-collector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.local.iterm2-gpu-collector</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>__INSTALL_DIR__/scripts/gpu_collector.py</string>
    <string>--cache-path</string>
    <string>__HOME__/.cache/iterm2-gpu/gpu_usage.json</string>
    <string>--interval</string>
    <string>2.5</string>
    <string>--allow-unavailable-write</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>__HOME__/.cache/iterm2-gpu/collector.out.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/.cache/iterm2-gpu/collector.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
